# systemd/mnt-unraid-vault.automount
[Unit]
Description=Automount Unraid NFS Share (Vault)
After=network-online.target
Wants=network-online.target

[Automount]
Where=/mnt/unraid/vault
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
