# systemd/mnt-unraid-shared.automount
[Unit]
Description=Automount Unraid NFS Share (Shared)
After=network-online.target
Wants=network-online.target

[Automount]
Where=/mnt/unraid/shared
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
