# systemd/mnt-unraid-foundryvtt.automount
[Unit]
Description=Automount Unraid NFS Share (FoundryVTT)
After=network-online.target
Wants=network-online.target

[Automount]
Where=/mnt/unraid/foundryvtt
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
