# systemd/mnt-unraid-data.automount
[Unit]
Description=Automount Unraid NFS Share (Data)
After=network-online.target
Wants=network-online.target

[Automount]
Where=/mnt/unraid/data
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
